# runtime/nas-systemd/sn9c292b-osd-off.service
[Unit]
Description=SN9C292B OSD off hotplug helper
Documentation=https://go2rtc.org/

[Service]
Type=oneshot
ExecStart=/opt/sn9c292b/osd-off.sh
